# scripts/pingan/weed-volume-01.ini
[program:weed-volume-01]
directory = /data1/www/dfs/volume01
command = /data1/www/weed/weed volume -ip=10.59.75.71 -port=18081 -mserver="10.59.75.71:9333" -dir=/data1/www/dfs/volume01 -dataCenter=dc1 -rack=rack1 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume01-info.log