# scripts/pingan/weed-volume-02.ini
[program:weed-volume-02]
directory = /data1/www/dfs/volume02
command = /data1/www/weed/weed volume -ip=10.59.75.71 -port=18082 -mserver="10.59.75.71:9333" -dir=/data1/www/dfs/volume02 -dataCenter=dc2 -rack=rack2 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume02-info.log