# scripts/pingan/whalefs-filer-01.ini
[program:whalefs-filer-01]
directory = /data1/www/whalefs
command = /data1/www/whalefs/whalefs filer
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
environment=ENV=development
stdout_logfile = /data1/log/dfs/whalefs-info.log