# scripts/prd/weed-volume-04.ini
[program:weed-volume-04]
directory = /data/dfs/volume04
command = /opt/weed/weed volume -ip=192.168.220.65 -port=18084 -mserver="192.168.220.65:9333,192.168.220.66:9333,192.168.220.65:9334" -dir=/data/dfs/volume04  -dataCenter=dc1 -rack=rack4  -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/weed/volume04-info.log