# scripts/prd/weed-master-01.ini
[program:weed-master-01]
directory = /data/dfs/master01
command = /opt/weed/weed master -ip=192.168.220.65 -port=9333 -peers="192.168.220.65:9333,192.168.220.66:9333,192.168.220.65:9334"  -mdir=/data/dfs/master01 -defaultReplication=100
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/weed/master01-info.log