# scripts/prd/whalefs-filer-01.ini
[program:whalefs-filer-01]
directory = /opt/framework/filer/current/
command = /opt/framework/filer/current/whalefs filer
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
environment=ENV=prd
stdout_logfile = /opt/framework/filer/current/info.log