# scripts/qa/weed-volume-06.ini
[program:weed-volume-06]
directory = /opt/dfs/volume06
command = /opt/weed/weed volume -ip=192.168.1.8 -port=18082 -mserver="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -dir=/opt/dfs/volume06 -dataCenter=dc1 -rack=rack2 -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/volume06/info.log