# scripts/qa/weed-volume-08.ini
[program:weed-volume-08]
directory = /opt/dfs/volume08
command = /opt/weed/weed volume -ip=192.168.1.8 -port=18084 -mserver="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -dir=/opt/dfs/volume08 -dataCenter=dc2 -rack=rack2 -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/volume08/info.log