# scripts/qa/whalefs-filer-01.ini
[program:whalefs-filer-01]
directory = /data1/www/filer
command = /data1/www/filer/whalefs filer
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
environment=ENV=prd
stdout_logfile = /data1/www/filer/info.log