# scripts/weed-volume-03.ini
[program:weed-volume-03]
directory = /opt/dfs/volume03
command = /opt/weed/weed volume -ip=192.168.1.9 -port=18083 -mserver="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -dir=/opt/dfs/volume03  -dataCenter=dc2 -rack=rack1  -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/volume03/info.log