# scripts/weed-volume-04.ini
[program:weed-volume-04]
directory = /opt/dfs/volume04
command = /opt/weed/weed volume -ip=192.168.1.9 -port=18084 -mserver="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -dir=/opt/dfs/volume04  -dataCenter=dc2 -rack=rack2  -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/volume04/info.log