# scripts/weed-volume-07.ini
[program:weed-volume-07]
directory = /opt/dfs/volume07
command = /opt/weed/weed volume -ip=192.168.1.8 -port=18083 -mserver="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -dir=/opt/dfs/volume07 -dataCenter=dc2 -rack=rack1 -max=30
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/volume07/info.log