# scripts/weed-master-03.ini
[program:weed-master-03]
directory = /opt/dfs/master03
command = /opt/weed/weed master -ip=192.168.1.9 -port=9335 -peers="192.168.1.9:9333,192.168.1.8:9333,192.168.1.9:9335" -mdir=/opt/dfs/master03  -defaultReplication=100
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/dfs/master03/info.log