# scripts/whalefs-task.ini
[program:whalefs-task]
directory = /opt/framework/task/current/
command = sudo /opt/framework/task/current/whalefs cron
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/framework/task/current/logs/info.log