# scripts/whalefs.ini
[program:mro-item-job]
directory = /opt/mro/mro-item-job/current/
command = sudo /opt/mro/mro-item-job/current/mro-item-job-exec.jar run --spring.profiles.active=pre
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = jenkins
stdout_logfile = /opt/mro/mro-item-job/logs/info.log