# scripts/whalefs-filer.ini
[program:whalefs-filer]
directory = /opt/framework/filer/current/
command = sudo /opt/framework/filer/current/whalefs filer
autostart = true
startsecs = 30
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user = root
stdout_logfile = /opt/framework/filer/current/logs/info.log