# scripts/pingan-prd/weed-volume-06.ini
[program:weed-volume-06]
directory = /data1/www/dfs/volume06
command = /data1/www/weed/weed volume -ip=zhc-zc-dfs002 -port=18082 -mserver="shc-zc-dfs001:9333,shc-zc-dfs002:9333,shc-zc-dfs003:9333" -dir=/data1/www/dfs/volume06 -dataCenter=dc2 -rack=rack22 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume06-info.log