# scripts/pingan-prd/weed-volume-07.ini
[program:weed-volume-07]
directory = /data1/www/dfs/volume07
command = /data1/www/weed/weed volume -ip=zhc-zc-dfs002 -port=18083 -mserver="shc-zc-dfs001:9333,shc-zc-dfs002:9333,shc-zc-dfs003:9333" -dir=/data1/www/dfs/volume07 -dataCenter=dc2 -rack=rack23 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume07-info.log