# scripts/pingan-prd/weed-volume-10.ini
[program:weed-volume-10]
directory = /data1/www/dfs/volume10
command = /data1/www/weed/weed volume -ip=zhc-zc-dfs003 -port=18082 -mserver="shc-zc-dfs001:9333,shc-zc-dfs002:9333,shc-zc-dfs003:9333" -dir=/data1/www/dfs/volume10 -dataCenter=dc3 -rack=rack32 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume10-info.log