# scripts/pingan-prd/weed-volume-04.ini
[program:weed-volume-04]
directory = /data1/www/dfs/volume04
command = /data1/www/weed/weed volume -ip=zhc-zc-dfs001 -port=18084 -mserver="shc-zc-dfs001:9333,shc-zc-dfs002:9333,shc-zc-dfs003:9333" -dir=/data1/www/dfs/volume04 -dataCenter=dc1 -rack=rack21 -max=30
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
stdout_logfile = /data1/log/dfs/volume04-info.log