# scripts/pingan-prd/whalefs-filer-04.ini
[program:whalefs-filer-04]
directory = /data1/www/whalefs04
command = /data1/www/whalefs04/whalefs filer
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
environment=ENV=ga
stdout_logfile = /data1/log/dfs/whalefs04-info.log