# scripts/pingan-prd/whalefs-filer-03.ini
[program:whalefs-filer-03]
directory = /data1/www/whalefs
command = /data1/www/whalefs/whalefs filer
autostart = true
startsecs = 10
autorestart = true
stopasgroup = true
stopsignal = KILL
killasgroup = true
user=www
environment=ENV=development
stdout_logfile = /data1/log/dfs/whalefs03-info.log